# roles/panda-service/templates/panda-service.conf
description "Panda Service"

start on runlevel [2345]
stop on runlevel [!2345]

exec {{ service_dir }}/panda-service/panda-service.py >> {{ log_file }} 2>&1